# oak_containers_system_image/files/etc/systemd/system/root-passwd.service
[Unit]
Description=Reset root password
#ConditionKernelCommandLine=debug

[Install]
WantedBy=getty.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo root:root|/usr/sbin/chpasswd'
